# systemd/user/stremio-vpn-watch.service
[Unit]
Description=Guard Stremio behind NordVPN
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/projects/streamio
ExecStart=%h/.local/bin/uv --cache-dir %h/projects/streamio/.uv-cache run python %h/projects/streamio/bin/stremio-vpn watchdog
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
